# naemon.ini
[supervisord]
nodaemon=true

[program:exim]
command=/usr/sbin/exim -bdf -q30m

[program:naemon]
command=/bin/sh -c "/usr/bin/naemon --verify-config /etc/naemon/naemon.cfg && /usr/bin/naemon /etc/naemon/naemon.cfg"
user=naemon
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
